# roles/borg_client/templates/borgbk.service.j2
# Generated by Ansible.
# This unit file is adapted from Borgmatic's unit file [1].
#
# [1]: https://projects.torsion.org/borgmatic-collective/borgmatic/raw/branch/main/sample/systemd/borgmatic.service)

[Unit]
Description=Borg Backup
Wants=network-online.target
After=network-online.target
ConditionACPower=true

[Service]
Type=oneshot

# for good luck...
# ExecStartPre=sleep 3m

ExecStart=systemd-inhibit --who="borgbackup" --what="sleep:shutdown" --why="Backup is running" {{ borg_client_user.home }}/borgbackup.sh

# Sandboxing settings
LockPersonality=true
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
# if you don't want to immediately kill, use EPERM here
SystemCallErrorNumber=kill
CapabilityBoundingSet=CAP_DAC_READ_SEARCH

# File paths
ProtectSystem=strict
BindPaths=-{{ borg_client_user.home }}/.cache/borg -{{ borg_client_user.home }}/.config/borg
BindReadOnlyPaths=-{{ borg_client_user.home }}/.ssh

# Low CPU priority
Nice=19

Restart=no
LogRateLimitIntervalSec=0

