# roles/borg_client/templates/borgbk.timer.j2
[Unit]
Description=Run Borg backup

[Timer]
# daily refers to midnigiht
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3h

[Install]
WantedBy=timers.target
